# scripts/scanoss-hfh-api.service
[Unit]
Description=SCANOSS HFH API Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=scanoss
ExecStart=/usr/local/bin/scanoss-hfh-api.sh

[Install]
WantedBy=multi-user.target